# roles/kube-worker/templates/kube-proxy.service.j2
[Unit]
Description=Kube proxy for Kubernetes

[Service]
ExecStart=/usr/bin/hyperkube kube-proxy --master=http://{{ master_ip }}
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target

